# home/ayamir/.config/autostart/qq.desktop
[Desktop Entry]
Name=QQ
Exec=/usr/bin/qq %U
Terminal=false
Type=Application
Icon=qq
StartupWMClass=QQ
Categories=Network;
Comment=QQ
